# packages/zbx.conf
[program:zabbix]
command=/usr/sbin/zabbix_server
autostart=true
autorestart=true
